# scripts-and-conf/etc/systemd/system/libvirt-guests.service.d/override.conf
# Placed in /etc/systemd/system/libvirt-guests.service.d/override.conf

[Service]
ExecStop=
ExecStop=/usr/local/bin/shutdown-vmx
ExecStop=/usr/lib/libvirt/libvirt-guests.sh stop
